# freeradius_configurator.upstart
description "Authentications"

start on starting access-control
stop  on stopped  nginx

env PYTHONPATH=/origin/lib/python

respawn
respawn limit 10 5

exec /origin/bin/authentication_freeradius_configurator


pre-start script
  rm -f /etc/freeradius/sites-enabled/default
  if ! grep -q '/origin/authentication/freeradius/dictionary' /etc/freeradius/dictionary
  then
    echo '' >> /etc/freeradius/dictionary
    echo '$INCLUDE /origin/authentication/freeradius/dictionary' >> /etc/freeradius/dictionary
  fi
end script
