# freeradius_configurator.upstart
description "Captive Portal"

start on starting access-control
stop  on stopped  nginx

env PYTHONPATH=/origin/lib/python

respawn
respawn limit 10 5

exec /origin/bin/captive-portal_freeradius_configurator


pre-start script
  if ! grep -q 'CP-Auth-ID' /etc/freeradius/dictionary
  then
    cat /origin/captive-portal/freeradius/dictionary >> /etc/freeradius/dictionary
  fi
end script