# rule-fetcher.cron
24 *	* * *	root	PYTHONPATH=/opt/origin/lib/python /opt/origin/sbin/rule-fetcher
