# debian/ea-authentication.freeradius.upstart
# freeradius -- Radius server
#

description     "Extensible, configurable radius daemon"
author          "Michael Vogt <mvo@ubuntu.com>"

start on started axon

respawn

env PYTHONPATH=/origin/lib/python

script
  if [ -r /etc/default/freeradius ]; then
    . /etc/default/freeradius
  fi
  exec /usr/sbin/freeradius -f $FREERADIUS_OPTIONS
end script

pre-start script
  # /var/run may be a tmpfs
  if [ ! -d /var/run/freeradius ]; then
    mkdir -p /var/run/freeradius
    chown freerad:freerad /var/run/freeradius
  fi

  rm -f /etc/freeradius/sites-enabled/default /etc/freeradius/sites-enabled/inner-tunnel
  perl -p -i -e "s/allow_vulnerable_openssl *= *no/allow_vulnerable_openssl = 'CVE-2014-0160'/" /etc/freeradius/radiusd.conf
  if ! grep -q '/origin/authentication/freeradius/dictionary' /etc/freeradius/dictionary
  then
    echo '' >> /etc/freeradius/dictionary
    echo '$INCLUDE /origin/authentication/freeradius/dictionary' >> /etc/freeradius/dictionary
  fi
end script
