# debian/ea-authentication.authentication_provider.upstart
description "Authentications"

start on starting access-control
stop  on stopped  nginx

env PYTHONPATH=/origin/lib/python

respawn
respawn limit 10 5

exec /origin/bin/authentication_provider


pre-start script
  rm -f /etc/freeradius/sites-enabled/default /etc/freeradius/sites-enabled/inner-tunnel
  if ! grep -q '/origin/authentication/freeradius/dictionary' /etc/freeradius/dictionary
  then
    echo '' >> /etc/freeradius/dictionary
    echo '$INCLUDE /origin/authentication/freeradius/dictionary' >> /etc/freeradius/dictionary
  fi
end script
