# debian/ea-core.nginx.upstart
description "nginx http daemon"
 
start on started redis
stop on runlevel [!2345]

env PYTHONPATH=/origin/lib/python
 
expect fork
respawn
respawn limit 10 5
 
exec /usr/sbin/nginx

pre-start script
  python -c 'from origin.neuron import Axon; Axon.generate_nginx_conf()'
end script
