# debian/ea-captive-portal.captive-portal_configuration_cacher.upstart
description "Captive Portal Conf retriever"

start on started axon

env PYTHONPATH=/origin/lib/python

respawn
respawn limit 10 5

exec /origin/bin/captive-portal_configuration_cacher

