# debian/ea-core.redis.upstart
description "redis server"

start on runlevel [2345]
stop on shutdown

respawn

setuid redis
setgid redis

exec /usr/bin/redis-server /origin/core/redis/conf

env PYTHONPATH=/origin/lib/python

post-start script
  python3 -c 'from origin.neuron import wait_for_synapse_ready; wait_for_synapse_ready()'
end script
