# debian/ea-core.nginx.upstart
description "nginx http daemon"
 
start on started redis
stop on runlevel [!2345]

env PYTHONPATH=/origin/lib/python
 
respawn
respawn limit 10 5
 
exec /usr/sbin/nginx -g 'daemon off;'

pre-start script
  rm -f /etc/nginx/sites-enabled/default
  python3 -c 'from origin.neuron import Axon; Axon.generate_nginx_conf()'
end script
