# debian/ea-core.axon.upstart
description "Axon"
 
start on started redis

respawn
respawn limit 10 5

env PYTHONPATH=/origin/lib/python

exec /origin/bin/axon

post-start script
  python3 -c 'from origin.neuron import wait_registered; wait_registered()'
end script

post-stop script
  # Make sure to mark WS disconnected....
  python3 -c 'from origin.neuron import Axon; Axon.redis_set_disconnected()'
end script
