# debian/ea-ids.rule-fetcher.cron.d
24 *	* * *	root	PYTHONPATH=/origin/lib/python /origin/sbin/rule-fetcher
