# debian/ea-ids.suricata.upstart
description	"Suricata server"

start on started axon

respawn
respawn limit 10 5
umask 022

reload signal USR2

exec /usr/bin/suricata --nflog


env PYTHONPATH=/origin/lib/python
pre-start script
    mkdir -p /var/log/suricata
    mkdir -p /etc/suricata/rules
    python3 -c 'from origin.ids import generate_suricata_conf; generate_suricata_conf(force=True)'
    /origin/bin/rule-fetcher --no-rule-reload || true
end script
