# debian/ea-ids.suricata.upstart
description	"Suricata server"

start on runlevel [2345]
stop on shutdown

respawn
respawn limit 10 5
umask 022

exec /usr/bin/suricata --pcap

pre-start script
    mkdir -p /var/log/suricata
    mkdir -p /etc/suricata/rules
    PYTHONPATH=/origin/lib/python /origin/sbin/rule-fetcher --no-rule-reload || true
end script
