# debian/ea-nac.snmptrapd.upstart
description "Poll Switches"

start on runlevel [2345]

respawn
respawn limit 10 5

exec /usr/sbin/snmptrapd -f -On -Lf /var/run/trap_fifo udp:162 udp6:162

pre-start script
  mkfifo /var/run/trap_fifo || true
end script