# debian/ea-nac.snmp_poller.upstart
description "Poll Switches"

start on started authentication_freeradius_configurator

env PERL5LIB=/origin/lib/perl5

respawn
respawn limit 10 5

exec /origin/bin/snmp_poller
