# debian/ea-nac.nac_configurator.upstart
description "Authentications"

start on started authentication_freeradius_configurator

env PYTHONPATH=/origin/lib/python

respawn
respawn limit 10 5

exec /origin/bin/nac_configurator
