# debian/elan-agent.session-trackerd.service
[Unit]
Description=Device Session Tracker

After=registered.service
Wants=registered.service

[Service]
Restart=always
RestartSec=10

Environment=PYTHONPATH=/origin/lib/python

ExecStart=/origin/bin/session-trackerd

[Install]
WantedBy=multi-user.target

