# debian/elan-agent.elan-radius-configurator.service
[Unit]

Description=RADIUS configurator


After=redis.service
Wants=redis.service

[Service]
Restart=always
RestartSec=10

Environment=PYTHONPATH=/elan-agent/lib/python

ExecStart=/elan-agent/bin/radius-configurator

[Install]
WantedBy=multi-user.target

