# debian/elan-agent.authentication_provider.upstart
description "Authentications"

start on started axon

env PYTHONPATH=/origin/lib/python

respawn
respawn limit 10 5

exec /origin/bin/authentication_provider


