# debian/elan-agent.rule-fetcher.cron.d
24 *	* * *	root	PYTHONPATH=/origin/lib/python /origin/bin/rule-fetcher
