# debian/elan-agent.mac_authz_manager.upstart
description "Allow MACs on VLANs"

start on started access-control

env PYTHONPATH=/origin/lib/python

respawn
respawn limit 10 5

exec /origin/bin/mac_authz_manager
