# debian/elan-agent.ids-monitor.service
[Unit]
Description=IDS monitor

[Service]
Restart=always
RestartSec=10

Environment=PYTHONPATH=/origin/lib/python

ExecStart=/origin/bin/ids-monitor

[Install]
WantedBy=multi-user.target
