# debian/postinst.in
# get rid of rc.d symlinks for freeradius
update-rc.d -f freeradius remove
rm -f /etc/init.d/freeradius
ln -s /lib/init/upstart-job /etc/init.d/freeradius

# create default params
/etc/freeradius/certs/bootstrap