# debian/postinst
# get rid of rc.d symlinks for freeradius
update-rc.d -f freeradius remove